# profile_files/dell-xps-9350-intel/system_files/usr/lib/systemd/user/purplefin-dell-xps-9350-panel.service
[Unit]
Description=Apply Dell XPS 13 9350 AC/battery panel policy
Documentation=file:/usr/share/purplefin/dell-xps-9350-panel.conf
After=graphical-session-pre.target
PartOf=graphical-session.target
ConditionPathExists=/sys/class/dmi/id/product_name

[Service]
Type=simple
ExecStart=/usr/libexec/purplefin/dell-xps-9350-panel-policy --watch
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=graphical-session.target
