# profile_files/dell-xps-9350-intel/system_files/usr/lib/systemd/system/purplefin-dell-ipu7-camera.service
[Unit]
Description=Activate Dell XPS 9350 IPU7 camera
After=systemd-modules-load.service
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/purplefin/dell-ipu7-activate
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
