# profile_files/dell-xps-9350-intel/system_files/usr/lib/systemd/system/purplefin-dell-ipu7-v4l2loopback-load.service
[Unit]
Description=Load Dell IPU7 v4l2loopback fallback device
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service
ConditionPathExists=/var/lib/purplefin/dell-ipu7/setup.done

[Service]
Type=oneshot
ExecStart=/usr/sbin/modprobe v4l2loopback
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
