# profile_files/dell-xps-9350-intel/system_files/usr/lib/systemd/system/purplefin-plymouth-initramfs.service
[Unit]
Description=Track Dell XPS Plymouth files in initramfs
After=rpm-ostreed.service
ConditionPathExists=/etc/plymouth/plymouthd.conf
ConditionPathExists=!/var/lib/purplefin/plymouth-initramfs-tracked

[Service]
Type=oneshot
ExecStart=/usr/libexec/purplefin/track-plymouth-initramfs

[Install]
WantedBy=multi-user.target
