# profile_files/dell-xps-9350-intel/system_files/usr/lib/systemd/user/pipewire.service.d/10-purplefin-dell-ipu7-libcamera.conf
[Service]
Environment=LD_LIBRARY_PATH=/usr/local/lib64:/usr/local/lib
Environment=LIBCAMERA_PLUGIN_PATH=/usr/local/lib64/libcamera:/usr/local/lib/libcamera
Environment=LIBCAMERA_IPA_MODULE_PATH=/usr/local/lib64/libcamera:/usr/local/lib/libcamera
Environment=LIBCAMERA_LOG_LEVELS=IPU7:0
Environment=LIBCAMERA_LOG_FILE=/tmp/libcamera-pw.log
