# profile_files/components/devops/system_files/usr/lib/systemd/user/purplefin-zsh-defaults.service
[Unit]
Description=Install Purplefin zsh defaults

[Service]
Type=oneshot
ExecStart=/usr/libexec/purplefin/install-zsh-defaults

[Install]
WantedBy=default.target
