# system_files/usr/lib/systemd/user/purplefin-ghostty-defaults.service
[Unit]
Description=Install Purplefin Ghostty defaults
ConditionPathExists=!%h/.config/ghostty/config.ghostty

[Service]
Type=oneshot
ExecStart=/usr/libexec/purplefin/install-ghostty-defaults

[Install]
WantedBy=default.target
